# aodv2.service
#
# Installed to /usr/lib/systemd/system/aodv2.service
#

[Unit]
Description=Linux Always-On Diagnostics Controller Service
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
# AOD_PYTHON selects the interpreter. On any host where the
# default python3 is too old, override this in /etc/aodv2/aodv2.env, e.g.
#   AOD_PYTHON=/usr/bin/python3.11        (RHEL 9 with python3.11 rpm)
#   AOD_PYTHON=/opt/aodv2-venv/bin/python (venv install)

Environment=AOD_PYTHON=/usr/bin/python3
EnvironmentFile=-/etc/aodv2/aodv2.env
ExecStart=/bin/sh -c 'exec "$AOD_PYTHON" -O /opt/aodv2/src/aod_entry.py'
Environment=PYTHONPATH=/opt/aodv2/src
Environment=AOD_CONFIG=/etc/aodv2/config.yaml
Restart=on-failure
RestartSec=5
User=root

KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=90s

OOMScoreAdjust=-500
LimitNOFILE=8192

[Install]
WantedBy=multi-user.target